# ubuntu_1804/ubuntu_1804.gitolite/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
command=/usr/sbin/sshd -D

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND

